# supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
nodaemon=false
logfile=/tmp/sandbox-logs/supervisord.log
pidfile=/tmp/supervisord.pid
user=claude

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:xvfb]
command=Xvfb :99 -screen 0 %(ENV_RESOLUTION)s -ac +extension GLX +render -noreset
priority=10
autorestart=true
startsecs=1
stdout_logfile=/tmp/sandbox-logs/xvfb.log
stderr_logfile=/tmp/sandbox-logs/xvfb.err

[program:x11vnc]
command=x11vnc -display :99 -forever -shared -nopw -rfbport 5900 -q
priority=20
autorestart=true
startsecs=2
stdout_logfile=/tmp/sandbox-logs/x11vnc.log
stderr_logfile=/tmp/sandbox-logs/x11vnc.err

[program:websockify]
command=websockify --web /usr/share/novnc 6080 localhost:5900
priority=30
autorestart=true
startsecs=2
stdout_logfile=/tmp/sandbox-logs/websockify.log
stderr_logfile=/tmp/sandbox-logs/websockify.err
